# package/install.sh
#!/bin/bash
chmod 755 /etc/init.d/boji
update-rc.d boji defaults